# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mcp-gateway]
command=pnpm start
directory=/app/mcp-gateway-registry
autostart=true
autorestart=true
stderr_logfile=/var/log/mcp-gateway.err.log
stdout_logfile=/var/log/mcp-gateway.out.log

[program:n8n-server]
command=pnpm start
directory=/app/n8n-mcp-server
autostart=true
autorestart=true
stderr_logfile=/var/log/n8n-server.err.log
stdout_logfile=/var/log/n8n-server.out.log

[program:freqtrade-server]
command=/app/freqtrade-mcp-server/.venv/bin/python -m freqtrade main
directory=/app/freqtrade-mcp-server
autostart=true
autorestart=true
stderr_logfile=/var/log/freqtrade-server.err.log
stdout_logfile=/var/log/freqtrade-server.out.log
